# systemd/system/start.service
[Unit]
Description=My script

[Service]
ExecStart=/usr/bin/start.sh

[Install]
WantedBy=multi-user.target 
