# systemd/system/getty@tty1.service.d/autologin.conf
[Service]
ExecStart=/usr/bin/gdm
ExecStart=-/sbin/agetty -o '-p -- \\u' --noclear --autologin temp --keep-baud 115200,57600,38400,9600 - $TERM
